# dot_config/powershell/profile.d/10-vars.ps1.tmpl
{{ if .os.macos -}}
# Homebrew on mac
Add-Content -Path $PROFILE.CurrentUserAllHosts -Value '$(/opt/homebrew/bin/brew shellenv) | Invoke-Expression'
{{- end }}
