# dot_config/powershell/profile.d/20-completions.ps1
# Pixi
if (Get-Command zoxide -errorAction SilentlyContinue) {
    (& pixi completion --shell powershell) | Out-String | Invoke-Expression
}
