# dot_config/powershell/profile.d/23-zoxide.ps1
# # Zoxide: https://github.com/ajeetdsouza/Zoxide
# if (Get-Command zoxide -errorAction SilentlyContinue) {
#   Invoke-Expression (& { (zoxide init powershell | Out-String) })
# }
