# dot_config/powershell/profile.d/25-get-childitems.ps1
Install-Module PSColor -Force
Import-Module PSColor
