# dot_config/powershell/profile.d/30-aliases.ps1
New-Alias -Name codei -Value code-insiders
